# rootfs/etc/services.d/node-2/run
#!/bin/bash

echo "* Starting node 2..."
/usr/bin/monerod --detach --testnet --p2p-bind-port 38080 --rpc-bind-port 38081 --no-igd --hide-my-port --data-dir /data/node-2 --p2p-bind-ip 127.0.0.1 --log-level 1 --add-exclusive-node 127.0.0.1:28080 --rpc-bind-ip 0.0.0.0 --confirm-external-bind --zmq-rpc-bind-port 38083

s6-svc -d /var/run/s6/services/node-2